# systemd/chromium-kiosk.service
[Unit]
Description=Chromium Kiosk - Queue Display
Documentation=https://github.com/yourusername/kiosk-checkpoint
After=network-online.target graphical.target kiosk.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
WorkingDirectory=/home/pi

# Start Chromium in kiosk mode
ExecStart=/usr/bin/chromium-browser \
    --noerrdialogs \
    --kiosk \
    --incognito \
    --disable-translate \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-component-extensions-with-background-pages \
    --disable-default-apps \
    --disable-extensions \
    --disable-web-resources \
    --no-first-run \
    --no-default-browser-check \
    --disable-breakpad \
    --disable-client-side-phishing-detection \
    --disable-component-update \
    --disable-popup-blocking \
    --disable-prompt-on-repost \
    --disable-sync \
    --metrics-recording-only \
    --mute-audio \
    http://localhost:8080

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chromium-kiosk

[Install]
WantedBy=graphical.target
