# scripts/serial_ws.service
[Unit]
Description=Serial to WebSocket Bridge for Kiosk
After=network-online.target

[Service]
Type=simple
User=moremeng
WorkingDirectory=/home/moremeng/kiosk-checkpoint
ExecStart=/usr/bin/python3 /home/moremeng/kiosk-checkpoint/scripts/serial_ws_bridge.py --device /dev/ttyUSB0 --baud 9600 --port 8765
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
